# recipes-core/initscripts/android-initscripts/android_sensorservice.upstart
description "Android Sensor Management Service"

start on stopped android_finish

respawn

script
	android-exec /system/bin/sensorservice
end script
