# recipes-core/initscripts/android-initscripts/android_mediaserver.upstart
description "Android Media Server"

start on stopped android_finish

respawn

script
	android-exec /system/bin/mediaserver
end script
