# recipes-core/initscripts/android-initscripts/android_drmserver.upstart
description "Android DRM Server"

start on stopped android_finish

respawn

script
	android-exec /system/bin/drmserver
end script
