# recipes-core/initscripts/android-initscripts/android_servicemanager.upstart
description "Android Service Manager"

start on stopped android_finish

respawn

script
	android-exec /system/bin/servicemanager
end script
