# recipes-webos/webos-initscripts/webos-initscripts/cpufreq-setup.upstart
start on stopped finish

script
	# We're trying to set the cpufreq governor for the first four cpu's if they are
	# available
	for n in `seq 0 3`; do
		if [ -e /sys/devices/system/cpu/cpu$n/cpufreq/scaling_governor ] ; then
			echo ondemand > /sys/devices/system/cpu/cpu$n/cpufreq/scaling_governor
		fi
	done
end script
