# recipes-multimedia/pulseaudio/pulseaudio/pulseaudio.upstart
description "PulseAudio"

start on stopped finish

exec pulseaudio --system --log-target=syslog
respawn
