# recipes-multimedia/pulseaudio/pulseaudio/pulseaudio.upstart
description "PulseAudio"

# The event is emitted by the android init system once the mediaserver and audioflinger
# are running and ready to be used.
start on android-finished

respawn

exec pulseaudio --log-target=syslog
