# recipes-connectivity/connman/connman/connman.upstart
description "ConnMan Network Connection Manager"

start on stopped finish

console output

# FIXME the list of interfaces to ignore should not be part of the initscript itself but
# stored in a configuration file in a machine specific way.
exec /usr/sbin/connmand -I usb0,rndis0,p2p0 -P ethernet -n
respawn
