# meta-luneos/recipes-core/dbus/dbus/dbus-session.service
[Unit]
Description=D-Bus Session Message Bus

[Service]
Type=forking
ExecStart=/usr/bin/startup-dbus-session.sh
