# meta-luneos/recipes-webos/maliit-framework-webos/files/maliit-server.path
[Unit]
Wants=maliit-server.service

[Path]
PathExists=/tmp/lsm-ready

[Install]
WantedBy=multi-user.target

