# recipes-connectivity/ofono/ofono/ofono.upstart
description "oFono - Open Source Telephony"

start on stopped finish

console output

exec /usr/sbin/ofonod -n
respawn
